# projects/omnichannel-chat-hub/deploy/systemd/omnichannel-api.service
[Unit]
Description=Omnichannel Chat Hub API
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/opc/projects/omnichannel-chat-hub
EnvironmentFile=/home/opc/projects/omnichannel-chat-hub/deploy/env/backend.env
ExecStart=/home/opc/projects/omnichannel-chat-hub/deploy/scripts/run-backend.sh
Restart=always
RestartSec=2
User=opc
Group=opc

[Install]
WantedBy=multi-user.target
